# app/crontab
#run python script every minutes. to make it run every 13 days, use : 0 0 13 * *
0 0 * * * python /app/cronscript.py > /proc/1/fd/1 2>/proc/1/fd/2
#You need empty line at the end